# config/crontab
# Обновление IP Cloudflare каждый день в 3:00
0 3 * * * root /usr/local/bin/update-cloudflare-ips.sh reload > /var/log/cloudflare-update.log 2>&1
